# schedules/com.claude.sprint-pulse.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.sprint-pulse</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>{{HOME}}/.local/bin/claude -p "/sprint-pulse --team ACE" --allowedTools "Bash,Read,Glob" --output-format text --no-session-persistence >> /tmp/claude-sprint-pulse.log 2>&amp;1</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>{{HOME}}/.claude/skills</string>

    <key>StandardOutPath</key>
    <string>/tmp/claude-sprint-pulse.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-sprint-pulse.log</string>
</dict>
</plist>
