# schedules/com.claude.bonusly-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.bonusly-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>{{HOME}}/.local/bin/claude -p "/bonusly-sync" --allowedTools "Bash,Read,Glob" --output-format text --no-session-persistence >> /tmp/claude-bonusly-sync.log 2>&amp;1</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Day</key><integer>4</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>{{HOME}}/.claude/skills</string>

    <key>StandardOutPath</key>
    <string>/tmp/claude-bonusly-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-bonusly-sync.log</string>
</dict>
</plist>
